# tools/uta-note-server.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.uta-note.local-server</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>npm</string>
    <string>run</string>
    <string>serve:local</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/takeshi/Documents/GitHub/開発/karaoke</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/uta-note-server.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/uta-note-server.log</string>
</dict>
</plist>
